# gitview/conf/cron/gitview-sync-projects.cron
0 0 * * 1-5 /usr/sbin/gitview-sync-projects
